# systemd/user/graphical.target
[Unit]
Description=GUI

[Install]
Alias=default.target
