# systemd/user/xorg.service
[Unit]
Description=Xorg server

[Service]
Type=simple
ExecStart=/usr/bin/startx

[Install]
WantedBy=graphical.target
